# sys-power/fake-battery-nut/files/fake-battery-nut.initd
#!/sbin/openrc-run
# Copyright 1999-2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Bridge a NUT-monitored UPS into a virtual battery for UPower"

command="/usr/bin/nut-to-fakebattery"
supervisor="supervise-daemon"
respawn_delay=5
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need upsd
	use upsmon
}
